# smoker.service
[Unit]
Description=BT Smoker Monitor
After=bluetooth.target network-online.target
Wants=bluetooth.target network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/bt-smoker-monitor
ExecStart=/home/pi/bt-smoker-monitor/venv/bin/python3 server.py --interval 30 --port 8888
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
